# deploy/systemd/neviim-api.service
[Unit]
Description=Neviim API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/apps/app-neviim/backend
EnvironmentFile=/home/ubuntu/apps/app-neviim/backend/.env
ExecStart=/home/ubuntu/apps/app-neviim/backend/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
